# cron.d/update_frame.sh
#!/bin/sh

curl http://localhost:9999/internal/update/frame